# minecraft_discord_bot/rootfs/etc/services.d/minecraft-discord-bot/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Minecraft Discord Bot..."

cd /app

exec uv run python app/serverbot.py
